# simpleshaper/files/simpleshaper.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 arteq

START=99

start() {
        /bin/ss-start
}
stop() {
        /bin/ss-stop
}
