# transmission/files/transmission.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 OpenWrt.org

START=80

start() {
	export TRANSMISSION_WEB_HOME=/www/transmission
	transmission-daemon -g /etc/transmission-daemon
}

stop() {
	pkill -9 transmission

